# deploy/leonard-bootstrap.service
[Unit]
Description=Leonard — fetch prod secrets from SSM to tmpfs
After=network-online.target
Wants=network-online.target
Before=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/leonard/scripts/fetch-prod-secrets.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
